# packages/qmp-system/files/etc/uci-defaults/qmp_banner
#!/bin/sh

cat /etc/qmp/qmp.banner >> /etc/banner
echo " ---------------------------------------------------" >> /etc/banner
echo " qMp $REL $COD (branch $BRANCH rev.$REV-$DATE)" >> /etc/banner
echo " ---------------------------------------------------" >> /etc/banner
echo " http://dev.qmp.cat/projects/qmp/repository/show?branch=$BRANCH&rev=$REV" >> /etc/banner
echo " ---------------------------------------------------" >> /etc/banner
