# packages/qmp-system/files/etc/uci-defaults/qmp_set_default_password
#!/bin/sh
echo -e "13f\n13f" | passwd root && exit 0
