# packages/qmp-system/files/etc/uci-defaults/qmp_wifi_health_cron_patch
#!/bin/sh

if ! cat /etc/crontabs/root | grep wifihealth; then
	echo "*/5 * * * * /usr/bin/wifihealth.sh >> /tmp/log/wifihealth.log" >> /etc/crontabs/root
	/etc/init.d/cron enable
	/etc/init.d/cron start
fi
