# packages/qmp-system/files/etc/uci-defaults/qmp_set_default_password
#!/bin/sh

# Set the default root password, "13f", if there is none and the router has not
# been configured by qMp yet.

FILE="/qmp_configured"

if [ ! -f $FILE ] && [ $(grep "root::0:0:99999:7:::" /etc/shadow) ]; then
  echo -e "13f\n13f" | passwd root && exit 0
else
  exit 0
fi

exit 1
